# deploy/tempmail.service
[Unit]
Description=Temporary inbox web and SMTP receiver
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=tempmail
Group=tempmail
WorkingDirectory=/opt/tempmail
EnvironmentFile=-/etc/tempmail/tempmail.env
Environment=TEMPMAIL_DB=/var/lib/tempmail/messages.db
Environment=TEMPMAIL_DOMAIN=mail.example.com
Environment=TEMPMAIL_TTL_HOURS=24
Environment=TEMPMAIL_HTTP_HOST=0.0.0.0
Environment=TEMPMAIL_HTTP_PORT=80
Environment=TEMPMAIL_SMTP_HOST=0.0.0.0
Environment=TEMPMAIL_SMTP_PORT=25
Environment=TEMPMAIL_MAX_MESSAGE_BYTES=5242880
Environment=TEMPMAIL_DAILY_SEND_LIMIT=5
Environment=TEMPMAIL_MAILBOX_DAILY_LIMIT=0
ExecStart=/usr/bin/python3 /opt/tempmail/tempmail.py
Restart=always
RestartSec=3
AmbientCapabilities=CAP_NET_BIND_SERVICE
CapabilityBoundingSet=CAP_NET_BIND_SERVICE
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=full
ReadWritePaths=/var/lib/tempmail

[Install]
WantedBy=multi-user.target
